# templates/alert_manager.sysv.erb
#!/bin/bash
#
#       /etc/rc.d/init.d/alert_manager
#
#       Daemonize the prometheus alert manager
#
# chkconfig:   2345 95 20
# description: The Alertmanager handles alerts sent by client applications such as the Prometheus server.
# processname: alert_manager
# pidfile: /var/run/alert_manager/pidfile

# Source function library.
. /etc/init.d/functions

DAEMON=<%= scope.lookupvar('prometheus::alert_manager::bin_dir') %>/alert_manager
PID_FILE=/var/run/alert_manager/alert_manager.pid
LOG_FILE=/var/log/alert_manager

[ -e /etc/sysconfig/alert_manager ] && . /etc/sysconfig/alert_manager

export GOMAXPROCS=${GOMAXPROCS:-2}

#
# Create the /var/run/alert_manager directory, which can live on a tmpfs
# filesystem and be destroyed between reboots.
#
mkrundir() {
        [ ! -d /var/run/alert_manager ] && mkdir -p /var/run/alert_manager
        chown <%= scope.lookupvar('prometheus::alert_manager::user') %> /var/run/alert_manager
}

#
# Create a PID file if it doesn't already exist, for clean upgrades
# from previous init-script controlled daemons.
#
KILLPROC_OPT="-p ${PID_FILE}"
mkpidfile() {
        # Create PID file if it didn't exist
        mkrundir
        [ ! -f $PID_FILE ] && pidofproc $DAEMON > $PID_FILE
        chown <%= scope.lookupvar('prometheus::alert_manager::user') %> /var/run/alert_manager
        if [ $? -ne 0 ] ; then
            rm $PID_FILE
            KILLPROC_OPT=""
        fi
}

start() {
        echo -n "Starting prometheus alert manager: "
        mkrundir
        [ -f $PID_FILE ] && rm $PID_FILE
        daemon --user=<%= scope.lookupvar('prometheus::alert_manager::user') %> \
            --pidfile="$PID_FILE" \
            "$DAEMON" -config.file=<%= scope.lookupvar('prometheus::alert_manager::config_file') %> -storage.path=<%= scope.lookupvar('prometheus::alert_manager::storage_path') %> <%= scope.lookupvar('prometheus::alert_manager::extra_options') %> >> "$LOG_FILE" &
        retcode=$?
        mkpidfile
        touch /var/lock/subsys/alert_manager
        return $retcode
}

stop() {
        DELAY=5 # seconds maximum to wait for a leave

        echo -n "Shutting down prometheus alert_manager: "
        mkpidfile

        alert_manager_pid=$(cat $PID_FILE)
        killproc $KILLPROC_OPT $DAEMON -INT
        retcode=$?

        # We'll wait if necessary to make sure the leave works, and return
        # early if we can.  If not, escalate to harsher signals.
        try=0
        while [ $try -lt $DELAY ]; do
        if ! checkpid $alert_manager_pid ; then
          rm -f /var/lock/subsys/alert_manager
          return $retcode
        fi
        sleep 1
          let try+=1
        done

        # If acting as a server, use a SIGTERM to avoid a leave.
        # This behavior is also configurable.  Avoid doing a "leave" because
        # having servers missing is a bad thing that we want to notice.
        #
        # A SIGTERM will mark the node as "failed" until it rejoins.
        # killproc with no arguments uses TERM, then escalates to KILL.
        killproc $KILLPROC_OPT $DAEMON
        retcode=$?

        rm -f /var/lock/subsys/alert_manager $PID_FILE
        return $retcode
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status -p ${PID_FILE} $DAEMON
        ;;
    restart)
        stop
        start
        ;;
    reload)
        mkpidfile
        killproc $KILLPROC_OPT $DAEMON -HUP
        ;;
    condrestart)
        [ -f /var/lock/subsys/alert_manager ] && restart || :
        ;;
    *)
        echo "Usage: alert_manager {start|stop|status|reload|restart}"
        exit 1
        ;;
esac
retcode=$?
# Don't let the [OK] get stomped on.
echo
exit $retcode
